# etc/supervisor/conf.d/program__cron.conf
[program:cron]
command=/usr/sbin/cron
redirect_stderr=true
stdout_logfile=/var/log/supervisor/program-cron.log
startsecs=0
