# etc/supervisor/conf.d/program__ocspd.conf
[program:ocspd]
autorestart=true
command=/usr/local/ocspd/run_ocspd.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/program-ocspd.log
user=ocspd
